# tests/images/zookeeper_jepsen/supervisor_zookeeper.conf
[program:zookeeper]
command=bash /var/lib/dist/zookeeper/start.sh
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
priority=5
redirect_stderr=true
stdout_logfile=/var/log/zookeeper/zookeeper.log
